# archiso/airootfs/etc/skel/.config/autostart/start_conky.sh.desktop
[Desktop Entry]
Exec=/home/aamir/.conky/start_conky.sh
Icon=dialog-scripts
Name=start_conky.sh
Path=
Type=Application
X-KDE-AutostartScript=true
